# vibrasense-edge.service
[Unit]
Description=VibraSense Edge Client
Documentation=https://github.com/vibrasense/rpi-edge-client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rpi-edge-client
Environment="PATH=/home/pi/rpi-edge-client/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/pi/rpi-edge-client/venv/bin/python3 /home/pi/rpi-edge-client/src/main.py /home/pi/rpi-edge-client/config

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vibrasense-edge

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/rpi-edge-client/data /home/pi/rpi-edge-client/logs /home/pi/rpi-edge-client/backups

[Install]
WantedBy=multi-user.target
